# channelbridge/config/crontab/elevenia_channel_cron
* * * * * php /app/channelbridge/elevenia/artisan schedule:run >> /dev/null 2>&1
