# channelbridge/config/supervisor/elevenia_channel_worker.ini
[program:elevenia_channel_worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/channelbridge/elevenia/artisan queue:work redis --sleep=3 --tries=3 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cb-elevenia-worker.log